# build/debian/ajime-agent.service
[Unit]
Description=Ajigent Edge Agent
Documentation=https://docs.ajime.io/agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStartPre=/usr/bin/pkill -x ajigent 2>/dev/null; true
ExecStart=/usr/local/bin/ajigent
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ajigent

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/etc/ajime /var/log/ajime
PrivateTmp=true

# Resource limits
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
